# systemd-update-images/update-nas4ipa.timer
[Unit]
Description=Auto-update NAS4IPA stable image

[Timer]
OnCalendar=Sat *-*-* 01:00:00
Persistent=true

[Install]
WantedBy=timers.target
