# systemd/uodate-bootx-x86_64.service
[Unit]
Description=Update bootc container images

[Service]
Type=oneshot
PIDFile=/run/user/%U/update-bootc-x86_64.pid
ExecStart=%h/bin/update-bootc
RemainAfterExit=yes
