# systemd/update-toolbox.timer
[Unit]
Description=Build updated toolbox image

[Timer]
OnCalendar=Sat *-*-* 09:00:00
Persistent=true

[Install]
WantedBy=timers.target
